# macos/com.llmdash.dashboard.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  llmdash LaunchAgent for macOS — the launchd equivalent of the Linux systemd unit.
  (The one-line installer, scripts/install-macos.sh, generates this for you.)

  Manual setup:
    1. which node   -> replace NODE_PATH below
    2. which codex  -> replace CODEX_PATH below
       IMPORTANT: CODEX_PATH must be an ABSOLUTE path (e.g. /Users/you/.local/bin/codex).
       launchd runs agents with a minimal PATH (/usr/bin:/bin:/usr/sbin:/sbin), so a
       bare "codex" can never resolve there and Codex limits would silently never load.
       (NODE_PATH must be absolute for the same reason.)
    3. replace PROJECT_DIR with your checkout path (e.g. /Users/you/llmdash)
    4. cp this file (drop the .example) to ~/Library/LaunchAgents/com.llmdash.dashboard.plist
    5. launchctl load -w ~/Library/LaunchAgents/com.llmdash.dashboard.plist
       (stop with: launchctl unload -w ~/Library/LaunchAgents/com.llmdash.dashboard.plist)

  To restrict strictly to your tailnet, add an LLMDASH_HOST entry set to your
  Tailscale IP (100.x.y.z) alongside LLMDASH_PORT below.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llmdash.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>NODE_PATH</string>
    <string>src/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>PROJECT_DIR</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LLMDASH_PORT</key>
    <string>8787</string>
    <key>LLMDASH_CODEX_CMD</key>
    <string>CODEX_PATH</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/llmdash.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/llmdash.err</string>
</dict>
</plist>
